# tracer_bringup/scripts/tracer.rules
ACTION=="add", SUBSYSTEM=="net", ENV{INTERFACE}=="can*", RUN+="/usr/bin/ip link set $name up type can bitrate 500000"